# luci/themes/luci-theme-ppip/root/etc/uci-defaults/30_luci-theme-ppip
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ppip=/luci-static/ppip
	set luci.main.mediaurlbase=/luci-static/ppip
	commit luci
EOF
exit 0
